# docker/supervisord-allinone.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/15/main -c config_file=/etc/postgresql/15/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql_err.log
priority=1

[program:redis]
command=/usr/bin/redis-server --bind 0.0.0.0 --port 6379 --maxmemory-policy allkeys-lru
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_err.log
priority=2

[program:qdrant]
command=/usr/local/bin/qdrant --config-path /app/docker/qdrant-config.yaml
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/qdrant.log
stderr_logfile=/var/log/supervisor/qdrant_err.log
priority=3
environment=QDRANT__SERVICE__HTTP_PORT="6333",QDRANT__LOG_LEVEL="INFO"

[program:backend]
command=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/backend
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log
priority=4
environment=DATABASE_URL="%(ENV_DATABASE_URL)s",REDIS_URL="%(ENV_REDIS_URL)s",QDRANT_URL="%(ENV_QDRANT_URL)s",CLAUDE_API_KEY="%(ENV_CLAUDE_API_KEY)s",VECTOR_SERVICE_TYPE="%(ENV_VECTOR_SERVICE_TYPE)s",ENVIRONMENT="%(ENV_ENVIRONMENT)s",PYTHONPATH="/app"

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
priority=5

